# setup/packages/saphir/saphir/boot/splash.initd
#!/sbin/openrc-run

name="splash"

depend() {
    before modules
}

start() {
    fbsplash -s /boot/Splash_1024_768.ppm -T 2
}

stop() {
    ebegin ""
    eend $?
}