# setup/packages/saphir/saphir-av-eset/av-controller.initd
#!/sbin/openrc-run

name="av-controller"
command_background="yes"

depend() {    
    after syslog
}

start() {
    ebegin "Starting Saphir ESET controller"

    start-stop-daemon --start --make-pidfile --exec /usr/bin/python3 --pidfile /var/run/$name.pid --background -- /usr/lib/saphir/bin/av-controller.py
}

stop() {
    ebegin "Stopping Saphir ESET controller"
    
    kill -9 $(cat /var/run/$name.pid)

    eend $?
}