# setup/packages/saphir/saphir-gui/gui.initd
#!/sbin/openrc-run

name="gui"
command_background="yes"
export DISPLAY=:0

depend() {    
    after syslog
}

start() {
    ebegin "Starting GUI"

    chmod +x /usr/lib/saphir/src/saphir/start.sh
    mkdir -p /var/log/saphir

    start-stop-daemon --start --make-pidfile --exec /usr/lib/saphir/src/saphir/start.sh --pidfile /var/run/$name.pid --background    
}

stop() {
    ebegin "Stopping GUI"
    
    kill -9 $(cat /var/run/$name.pid)

    eend $?
}